# deploy/com.noto.resync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Noto nightly corpus re-sync — picks up newly added candidates/docs.
     Install:
       cp deploy/com.noto.resync.plist ~/Library/LaunchAgents/
       launchctl load -w ~/Library/LaunchAgents/com.noto.resync.plist
     Runs daily at 03:30. Logs: lark/resync.log
     Edit paths below if the repo is not at /Users/noto/noto-home. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.noto.resync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/noto/noto-home/tools/noto-resync.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/noto/noto-home</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>LOLABOT_HOME</key>
    <string>/Users/noto/noto-home</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/noto/noto-home/lark/resync.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/noto/noto-home/lark/resync.log</string>
</dict>
</plist>
